# vocore2/files/etc/uci-defaults/50-vocore2-config
#!/bin/sh

. /lib/functions.sh
. /lib/functions/system.sh

uci set system.@system[-1].hostname="vocore2"
uci set network.wwan=interface
uci set network.wwan.proto=dhcp

uci commit
